# scripts/xvfb.service
[Unit]
Description=Xvfb Virtual Display for IBKR Gateway
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/Xvfb :99 -screen 0 1024x768x24
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
